# systemd/deep-agg-dom.timer
[Unit]
Description=Run Deep Agg DOM pipeline nightly

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
